# src/etc/systemd/system/xiva-audio-config.service
[Unit]
Description=xiva-audio-config
 
[Service]
ExecStart=/usr/bin/xiva-audio-config
 
[Install]
WantedBy=multi-user.target
